# deploy/launchd/com.zephyr-systems.mcp-memory-libravdb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zephyr-systems.mcp-memory-libravdb</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/mcp-memory-libravdb</string>
    <string>http</string>
    <string>--http-port</string>
    <string>8082</string>
    <string>--http-host</string>
    <string>127.0.0.1</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>LIBRAVDBD_ADDR</key>
    <string>unix:///opt/homebrew/var/libravdbd/run/libravdb.sock</string>
    <key>LIBRAVDB_TENANT_KEY</key>
    <string>libravdb-mcp-server</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/opt/homebrew/var/log/mcp-memory-libravdb.log</string>
  <key>StandardErrorPath</key>
  <string>/opt/homebrew/var/log/mcp-memory-libravdb.err</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
